# .config/systemd/user/agenda-notify.service
[Unit]
Description=Morning dunst notification with today's org-mode agenda

[Service]
Type=oneshot
ExecStart=%h/.config/sway/agenda.sh notify
