# .config/systemd/user/agenda-refresh.service
[Unit]
Description=Refresh cached org-mode agenda for swaybar

[Service]
Type=oneshot
ExecStart=%h/.config/sway/agenda.sh refresh
