# .config/systemd/user/battery-notify.service
[Unit]
Description=Warn via dunst when the battery is low

[Service]
Type=oneshot
ExecStart=%h/.config/sway/battery-notify.sh
