# .config/systemd/user/battery-notify.timer
[Unit]
Description=Periodic low-battery check

[Timer]
OnBootSec=2min
OnUnitActiveSec=2min
Unit=battery-notify.service

[Install]
WantedBy=timers.target
